# release/src/router/swrtmesh/libeasy/source/gitlab-ci/iopsys-supervisord.conf
[program:ubusd]
priority=1
autorestart=false
startretries=0
command=/bin/bash -c "/usr/sbin/ubusd"
